# supervisord.conf
[program:supercronic]
command=/usr/local/bin/supercronic /app/camcommander/cron-camcommander

[program:imageshepherd]
command=imageshepherd
